# scripts/restic/systemd/restic-verify.timer
[Unit]
Description=Monthly Restic repository integrity check
Requires=restic-verify.service

[Timer]
# Run monthly on the 1st at 5:00 AM
OnCalendar=*-*-01 05:00:00
Persistent=true

[Install]
WantedBy=timers.target
